# deploy/supervisor/mes-backend.conf
[program:mes-backend]
# Supervisor 托管 Java 进程：开机自启、异常退出自动拉起，并把标准输出/错误写入独立日志。
# 后端只监听 127.0.0.1:8080，关闭桌面浏览器启动；公网流量必须经过 Nginx。
command=/bin/bash -lc 'MES_HOST=127.0.0.1 MES_PORT=8080 MES_OPEN_BROWSER=false /www/wwwroot/mes/scripts/run-backend.sh'
# Maven 从 backend 工作目录读取模块 pom，run-backend.sh 会再定位根目录 .env。
directory=/www/wwwroot/mes/backend
user=root
# 服务随 Supervisor 启动，非正常退出时自动重启；连续存活 5 秒才视为启动成功。
autostart=true
autorestart=true
startsecs=5
stopwaitsecs=30
# 停止整个 Maven/Java 进程组，避免重启后遗留旧端口监听进程。
stopasgroup=true
killasgroup=true
# 应用日志与异常日志分开保存，便于服务器排错。
stdout_logfile=/www/wwwlogs/mes-backend-access.log
stderr_logfile=/www/wwwlogs/mes-backend-error.log
# 显式固定 JDK 21 和可执行路径，避免 Supervisor 的非登录环境找不到 java/mvn。
environment=JAVA_HOME="/usr/lib/jvm/java-21-openjdk-amd64",PATH="/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin"
