# deploy/supervisor/mes-backend.conf
[program:mes-backend]
command=/bin/bash -lc 'MES_HOST=127.0.0.1 MES_PORT=8080 MES_OPEN_BROWSER=false /www/wwwroot/mes/scripts/run-backend.sh'
directory=/www/wwwroot/mes/backend
user=root
autostart=true
autorestart=true
startsecs=5
stopwaitsecs=30
stopasgroup=true
killasgroup=true
stdout_logfile=/www/wwwlogs/mes-backend-access.log
stderr_logfile=/www/wwwlogs/mes-backend-error.log
environment=JAVA_HOME="/usr/lib/jvm/java-21-openjdk-amd64",PATH="/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin"
